# board/common/overlay/etc/crontabs/root
13 * * * * /usr/sbin/logrotate /etc/logrotate.conf
0 2 * * 0 /usr/sbin/https-update
*/15 * * * * /usr/sbin/dyndns-update
*/15 * * * * date +%s > /var/lib/saved_date
*/5 * * * * (sleep $(expr $RANDOM % 300); /sbin/fwupdate-check) >> /var/log/fwupdate-check.log
* * * * * /sbin/send-sys-info &>> /var/log/send-sys-info.log
14 3 * * * /usr/libexec/reboot-when-idle
35 * * * * grep 'PERIODIC_RESET_PEERS=true' /data/etc/miner.conf && /sbin/resetpeers
